# flask_demo/flask_supervisord.conf
[supervisord]
nodaemon=true

[program:helloword]
command=python myproject.py
